# deploy/apk-landing.service
[Unit]
Description=apk-landing 京东落地页服务
After=network.target

[Service]
Type=simple
# 部署目录，按实际路径修改
WorkingDirectory=/opt/apk-landing
ExecStart=/opt/apk-landing/apk-landing -f /opt/apk-landing/etc/config.yaml
Restart=always
RestartSec=3
# 以普通用户运行更安全，按需修改；监听 5001 非特权端口无需 root
User=www-data
Group=www-data
# 日志走 journald，用 journalctl -u apk-landing 查看
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
