# deploy/apk-landing.service
[Unit]
Description=apk-landing 京东落地页服务
After=network.target

[Service]
Type=simple
# 部署目录
WorkingDirectory=/home/sysadmin/data/apk-landing/dist-linux
ExecStart=/home/sysadmin/data/apk-landing/dist-linux/apk-landing -f /home/sysadmin/data/apk-landing/dist-linux/etc/config.yaml
Restart=always
RestartSec=3
# 运行用户，监听 5001 非特权端口无需 root
User=sysadmin
Group=sysadmin
# 日志走 journald，用 journalctl -u apk-landing 查看
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
